# templates/etc/profile.d/java.sh.j2
export JAVA_HOME={{java_home}}
export JDK_HOME={{java_home}}
